# systemd/ssh-blocker.service.in
[Unit]
Description=Read logs and block attacks with ipset.
Before=iptables.service

[Service]
Type=simple
Restart=always
RestartSec=1s
ExecStart=@sbindir@/ssh-blocker ssh-blocker

# Read _IPv4_ address from /etc/ssh/sshd_config in the form:
# AllowUsers foo@1.2.3.4 bar@5.6.7.8
ExecStartPost=-/bin/sh -c "for ip in $(grep ^AllowUsers /etc/ssh/sshd_config | grep -Po '\w+@\K[\d.]+'); \
do /usr/sbin/ipset -exist add ssh-whitelist $ip timeout 0; done"

[Install]
WantedBy=multi-user.target
